# systemd/user/clipmenud.service.d/override.conf
[Service]
Environment=CM_DIR=%h/.cache/
